# roles/django/files/django_custom.service.j2
[Unit]
Description={{ project_name }}
After=postgresql.service

[Service]
TimeoutStartSec=0
User={{ django_user }}
WorkingDirectory={{ django_project_root }}
ExecStart={{ django_venv_dir }}/bin/python {{ item.command }}
PIDFile={{ django_project_root }}/{{ project_name }}_{{ item.name }}.pid

[Install]
WantedBy=multi-user.target